# quadlet/lme-dashboard.container
# lme-dashboard.container
[Unit]
Description=LME Security Dashboard - FastAPI + Alpine.js UI
Requires=lme-elasticsearch.service lme-litellm.service lme-pgvector.service lme-embeddings.service
After=lme-elasticsearch.service lme-litellm.service lme-pgvector.service lme-embeddings.service
PartOf=lme.service

[Install]
WantedBy=default.target lme.service

[Service]
Restart=always
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=10
TimeoutStartSec=5400
Environment=ANSIBLE_VAULT_PASSWORD_FILE=/etc/lme/pass.sh

[Container]
ContainerName=lme-dashboard
Image=localhost/lme-dashboard:LME_LATEST
Network=lme
PodmanArgs=--network-alias lme-dashboard
PublishPort=8502:8502
Environment=ELASTICSEARCH_URL=https://lme-elasticsearch:9200
Environment=ELASTICSEARCH_USER=elastic
Environment=LITELLM_URL=https://lme-litellm:4000
Environment=LITELLM_API_KEY=sk-lme-llama-proxy
Environment=LITELLM_MODEL=lfm2.5-1.2b-instruct
Environment=PGVECTOR_HOST=lme-pgvector
Environment=PGVECTOR_PORT=5432
Environment=PGVECTOR_DB=lme_vectors
Environment=PGVECTOR_USER=lme
Environment=EMBED_URL=https://lme-embeddings:8081
Volume=/opt/lme/lme-dashboard/app.py:/app/app.py:ro
Volume=/opt/lme/lme-dashboard/static:/app/static:ro
Volume=/opt/lme/config:/opt/lme/config:rw
Volume=/opt/lme/llama-models:/opt/lme/llama-models:rw
Volume=/opt/lme/scripts/sigma:/opt/lme/scripts/sigma:rw
Volume=/etc/lme/pass.sh:/etc/lme/pass.sh:ro
Volume=lme_certs:/certs:ro
Environment=LLAMA_MODELS_DIR=/opt/lme/llama-models
Environment=LLAMA_MODEL_CONFIG=/opt/lme/config/llama-cpp-model.json
Environment=LLAMA_MODEL_TRIGGER=/opt/lme/config/.llama-model-updated
Environment=LLAMA_MODEL_STATUS=/opt/lme/config/llama-cpp-status.json
Environment=SSL_KEYFILE=/certs/dashboard/dashboard.key
Environment=SSL_CERTFILE=/certs/dashboard/dashboard.crt
Environment=LITELLM_CONFIG_PATH=/opt/lme/config/litellm_config.yaml
Environment=LLM_KEYS_PATH=/opt/lme/config/llm_keys.enc
Environment=LLM_KEYS_TRIGGER=/opt/lme/config/.llm-keys-updated
Environment=VAULT_PASS_FILE=/etc/lme/pass.sh
Exec=uvicorn app:app --host 0.0.0.0 --port 8502 --ssl-keyfile /certs/dashboard/dashboard.key --ssl-certfile /certs/dashboard/dashboard.crt
Secret=elastic,type=env,target=ELASTICSEARCH_PASSWORD
Secret=pgvector,type=env,target=PGVECTOR_PASS
