# systemd/lme-llm-keys.service
[Unit]
Description=Sync LLM API keys and restart LiteLLM
After=lme-litellm.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/lme/scripts/sync_llm_keys.py
Environment=LLM_KEYS_PATH=/opt/lme/config/llm_keys.enc
Environment=VAULT_PASS_FILE=/etc/lme/pass.sh
