# systemd/lme-llama-model.path
[Unit]
Description=Watch for llama.cpp model switch requests

[Path]
PathModified=/opt/lme/config/.llama-model-updated

[Install]
WantedBy=multi-user.target
