# systemd/lme-llm-keys.path
[Unit]
Description=Watch for LLM API key changes

[Path]
PathModified=/opt/lme/config/.llm-keys-updated

[Install]
WantedBy=multi-user.target
